# ansible/templates/app.service.j2
[Unit]
Description=CI/CD Demo Flask App
After=network.target

[Service]
Type=simple
WorkingDirectory={{ app_dir }}
Environment="PATH={{ venv_dir }}/bin"
ExecStart={{ venv_dir }}/bin/gunicorn -w 2 -b 0.0.0.0:5000 app.app:app
Restart=on-failure

[Install]
WantedBy=multi-user.target
